# root/etc/services.d/2inputGrabber/run
#!/usr/bin/with-contenv bash

exec /bin/bash /app/bbuddy/example/grabInput.sh
