# root-i386/etc/services.d/0redis/run
#!/usr/bin/with-contenv bash

exec /usr/bin/redis-server /etc/redis.conf --daemonize no
